# ansible/roles/jenkins-infrastructure/templates/jenkins-master.service.j2
[Unit]
Description=Jenkins Blue-Green Multi-Team Coordinator Service
Documentation=https://jenkins.io/doc/
{% if jenkins_container_runtime is defined and jenkins_container_runtime != '' %}
Requires={{ jenkins_container_runtime }}.service
After={{ jenkins_container_runtime }}.service network-online.target
PartOf={{ jenkins_container_runtime }}.service
{% else %}
After=network-online.target
{% endif %}
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=no
TimeoutStartSec=600
TimeoutStopSec=60

# Pre-start health checks and setup
ExecStartPre=/bin/bash -c 'mkdir -p {{ jenkins_home_dir }}/logs {{ jenkins_home_dir }}/scripts'
ExecStartPre=/bin/bash -c 'chown -R {{ jenkins_user }}:{{ jenkins_group }} {{ jenkins_home_dir }}'
ExecStartPre={{ jenkins_home_dir }}/bin/jenkins-prestart.sh

# Verify containers are ready before declaring service active
ExecStart=/bin/bash -c 'echo "Jenkins Blue-Green Multi-Team service started at $(date)" >> {{ jenkins_home_dir }}/logs/service.log'
ExecStart=/bin/bash -c 'echo "Verifying Jenkins team containers..." >> {{ jenkins_home_dir }}/logs/service.log'

# Health check all team environments (with retry logic)
ExecStartPost=/bin/bash -c 'sleep 15'
ExecStartPost=/bin/bash -c 'for i in {1..3}; do if {{ jenkins_home_dir }}/bin/jenkins-healthcheck.sh; then break; fi; echo "Health check attempt $i failed, retrying..." >> {{ jenkins_home_dir }}/logs/service.log; sleep 10; done'

# Service shutdown (containers managed by Ansible)
ExecStop=/bin/bash -c 'echo "Jenkins Blue-Green Multi-Team service stopped at $(date)" >> {{ jenkins_home_dir }}/logs/service.log'

# Security settings
User={{ jenkins_user }}
Group={{ jenkins_group }}
NoNewPrivileges=yes

# Environment
Environment="JENKINS_HOME={{ jenkins_home_dir }}"
Environment="JENKINS_USER={{ jenkins_user }}"
Environment="CONTAINER_RUNTIME={{ jenkins_container_runtime }}"

# Working directory
WorkingDirectory={{ jenkins_home_dir }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jenkins-master

[Install]
WantedBy=multi-user.target