# ansible/roles/podman-v2/templates/podman-service.service.j2
# Systemd service template for Podman containers
# Generated by Ansible - DO NOT EDIT MANUALLY
# This is an example service file for managing Podman containers with systemd

[Unit]
Description=Podman {{ podman_service_name | default('example') }} Service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor={{ podman_service_volume_path | default('/var/lib/containers/storage') }}

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart={{ podman_service_restart | default('on-failure') }}
RestartSec={{ podman_service_restart_sec | default('30') }}
TimeoutStartSec={{ podman_service_timeout_start | default('70') }}
TimeoutStopSec={{ podman_service_timeout_stop | default('70') }}

# Security settings
User={{ podman_service_user | default('root') }}
Group={{ podman_service_group | default('root') }}

# Resource limits
{% if podman_service_memory_limit is defined %}
MemoryLimit={{ podman_service_memory_limit }}
{% endif %}

{% if podman_service_cpu_limit is defined %}
CPUQuota={{ podman_service_cpu_limit }}
{% endif %}

# Execute commands
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=no-conmon \
    --rm \
    --sdnotify=conmon \
    --replace \
    {{ podman_service_run_args | default('') }} \
    {{ podman_service_image | default('registry.access.redhat.com/ubi8/ubi:latest') }}

ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm -f --ignore --cidfile=%t/%n.ctr-id

# Podman-specific settings
Type=notify
NotifyAccess=all
KillMode=mixed

# Additional security
PrivateTmp={{ podman_service_private_tmp | default('yes') }}
ProtectSystem={{ podman_service_protect_system | default('strict') }}
ProtectHome={{ podman_service_protect_home | default('yes') }}
NoNewPrivileges={{ podman_service_no_new_privileges | default('yes') }}

[Install]
WantedBy={{ podman_service_wanted_by | default('multi-user.target') }}