# ansible/roles/monitoring/templates/monitoring-agent-health.service.j2
[Unit]
Description=Monitoring Agent Health Check
Documentation=https://github.com/your-org/jenkins-ha
After=docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring/agent-health-check.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitoring-agent-health

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/monitoring-agents

[Install]
WantedBy=multi-user.target
