# ansible/roles/high-availability/templates/jenkins-haproxy.service.j2
[Unit]
Description=Jenkins HAProxy Container Service
Documentation=https://github.com/haproxy/haproxy
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=forking
Restart=always
RestartSec=30
TimeoutStartSec=0
TimeoutStopSec=120

# Working directory
WorkingDirectory=/etc/haproxy

# Pre-start checks
ExecStartPre=-/usr/bin/docker-compose down
ExecStartPre=/usr/bin/docker-compose pull --quiet

# Start container
ExecStart=/usr/bin/docker-compose up -d

# Health check
ExecStartPost=/bin/bash -c 'timeout 60 bash -c "until docker ps --format \"table {{.Names}}\" | grep -q jenkins-haproxy; do sleep 2; done"'

# Stop container
ExecStop=/usr/bin/docker-compose down

# Reload configuration
ExecReload=/usr/bin/docker-compose restart

# Log output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jenkins-haproxy

# Security settings
User=root
Group=haproxy

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=COMPOSE_PROJECT_NAME=jenkins-haproxy
Environment=COMPOSE_HTTP_TIMEOUT=300
Environment=DOCKER_CLIENT_TIMEOUT=300

[Install]
WantedBy=multi-user.target
WantedBy=docker.service